# flatpak/navetted.service
[Unit]
Description=Relay — remote Claude Code session gateway
After=network.target

[Service]
ExecStart=/app/bin/navetted
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
